# private/src/nebula-api.service
[Unit]
Description=Nebula API (Project Zomboid Management)
After=network-online.target
Wants=network-online.target

[Service]
User=root
WorkingDirectory=/opt/nebula-api
EnvironmentFile=/opt/nebula-api/.env
ExecStart=/opt/nebula-api/.venv/bin/python /opt/nebula-api/nebula_api.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
